# content/etc/s6-overlay/s6-rc.d/fans/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

/app/fan_speeds.sh
